# _pre_posts/AppArmor/new_build/ros2_security_report.timer
[Unit]
Description=Run ROS2 Security Report Hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target